# ponix4-src/source/clients/citrix/linuxx86/linuxx86.cor/usb/ica-usb.rules
# Don't allow HAL to handle devices that ICA client is remoting

SUBSYSTEM=="usb", \
ACTION=="add", \
PROGRAM="###ICAROOT###/ctx_usb_isactive", \
    OPTIONS+="last_rule", \
    RUN+="/usr/bin/logger -t ica-usb.rules suppressing HAL action for %k"
