# ponix4-const/packages/sane/etc/udev/rules.d/10-scanner_aliasing.rules
SUBSYSTEMS=="usb",ATTRS{idVendor}=="04c5",ATTRS{idProduct}=="10ef",RUN+="/etc/udev/scripts/scanner_fj5015c"
