# ponix4-const/tmp-tree/etc/inittab
# /etc/inittab init(8) configuration for BusyBox

tty1::sysinit:/etc/ponix.init
tty1::respawn:/bin/start-session 0
tty12::askfirst:/bin/getty tty12 9600

::shutdown:/bin/shutdown

::ctrlaltdel:/bin/reboot
