# ponix4-const/packages/sane/etc/xinetd.d/sane
service sane 
{
	disable     = no
	socket_type = stream
	wait        = no
	user        = root
	server      = /bin/saned
	protocol    = tcp
}
