# ponix4-src/source/plymouth-0.8.8/systemd-units/plymouth-quit-wait.service.in
[Unit]
Description=Wait for Plymouth Boot Screen to Quit
After=rc-local.service plymouth-start.service systemd-user-sessions.service

[Service]
ExecStart=-@PLYMOUTH_CLIENT_DIR@/plymouth --wait
Type=oneshot
TimeoutSec=20
