# ponix4-src/source/plymouth-0.8.8/systemd-units/plymouth-quit.service.in
[Unit]
Description=Terminate Plymouth Boot Screen
After=rc-local.service plymouth-start.service systemd-user-sessions.service

[Service]
ExecStart=-@PLYMOUTH_CLIENT_DIR@/plymouth quit
Type=oneshot
TimeoutSec=20
