# ponix4-src/source/plymouth-0.8.8/systemd-units/plymouth-reboot.service.in
[Unit]
Description=Show Plymouth Reboot Screen
After=getty@tty1.service prefdm.service plymouth-start.service
Before=reboot.service
DefaultDependencies=no
ConditionKernelCommandLine=!plymouth.enable=0

[Service]
ExecStart=@PLYMOUTH_DAEMON_DIR@/plymouthd --mode=shutdown --attach-to-session
ExecStartPost=-@PLYMOUTH_CLIENT_DIR@/plymouth show-splash
Type=forking
