# ponix4-src/source/scard/RuToken/ponix/etc/udev/rules.d/95-rutokens.rules
# udev rules for Rutoken S
SUBSYSTEM=="usb", ACTION=="add", ATTRS{idVendor}=="0a89", ATTRS{idProduct}=="0020", RUN+="/bin/chgrp pcscd $root/$parent"
