# ponix4-const/tmp-tree/etc/udev/rules.d/55-udev-monitors.rules

ACTION=="change", SUBSYSTEM=="drm", RUN+="/bin/apply_xrandr"
