# ponix4-const/tmp-tree/etc/udev/rules.d/70-start-actkbd.rules
ACTION=="add", SUBSYSTEM=="input", KERNEL=="event[0-9]*", RUN+="/bin/actkbd -D -q -d %E{DEVNAME}"